# static/common/create_config.service
[Service]
ExecStart=/usr/local/bin/create_config.sh
Type=oneshot

[Install]
WantedBy=multi-user.target
